# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
Type=simple
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin n0tr00t %I $TERM
